# ecnu-portal.service
[Unit]
Description=ECNU Portal

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ecnu-portal login

