# ecnu-portal.timer
[Unit]
Description=ECNU Portal

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target

